# packaging/linux/debian/postrm
#!/bin/sh
set -e
if [ "$1" = "purge" ]; then
  systemctl disable offs-daemon 2>/dev/null || true
  systemctl daemon-reload
fi
