# packaging/linux/debian/postinst
#!/bin/sh
set -e
systemctl daemon-reload
systemctl enable offs-daemon 2>/dev/null || true
systemctl start offs-daemon 2>/dev/null || true
